# etc/change_hardcoded_username_first/etc/init.d/mts-ENLIVEN.conf
start on (filesystem and stopped udevtrigger)
stop on runlevel [06]

setuid owner
respawn

exec /home/owner/mts-ENLIVEN
